# socks5-server.service
[Unit]
Description=Socks5 Server
[Service]
ExecStart=/usr/bin/socks5-server
Restart=always
User=root
Group=root
[Install]
WantedBy=multi-user.target
