# debian/postinst
#!/usr/bin/env bash

set -e

systemctl enable socks5-server --now || true