# debian/prerm
#!/usr/bin/env bash

set -e

systemctl disable socks5-server --now || true